# docker/8.3-prod/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=admin
password=admin123

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=admin
password=admin123

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/tmp/nginx-supervisor.log
stopwaitsecs=10
priority=100

[program:php-fpm]
command=/usr/sbin/php-fpm8.3 --nodaemonize --fpm-config /etc/php/8.3/fpm/php-fpm.conf
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/tmp/php-fpm-supervisor.log
stopwaitsecs=10
priority=200

[program:queue]
command=bash -c "cd /var/www/html && /usr/bin/php artisan queue:work --sleep=3 --tries=3 --max-time=3600"
autostart=true
autorestart=true
user=www-data
directory=/var/www/html
redirect_stderr=true
stdout_logfile=/tmp/queue-supervisor.log
stopwaitsecs=360
priority=300
environment=HOME="/var/www/html",USER="www-data"

[program:schedule]
command=bash -c "cd /var/www/html && /usr/bin/php artisan schedule:work"
autostart=true
autorestart=true
user=www-data
directory=/var/www/html
redirect_stderr=true
stdout_logfile=/tmp/schedule-supervisor.log
stopwaitsecs=10
priority=400
environment=HOME="/var/www/html",USER="www-data"
